# data/mapplauncherd.target
[Unit]
Description=Application launch boosters
