# data/mapplauncherd.target
[Unit]
Description=Application launch boosters
Requires=pre-user-session.target
After=pre-user-session.target
